# asucfunding/deployment/xlab_uwsgi.conf
# file: /etc/init/uwsgi.conf

description "uWSGI starter"

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn

exec /opt/venvs/xlab_env/bin/uwsgi --uid www-data --home /opt/venvs/xlab_env --pythonpath /opt/django-trunk/xlab_server --socket /tmp/uwsgi_xlab.sock --chmod-socket --module wsgi --logdate --processes 2 --master --logto /opt/log/uwsgi/xlab.log --optimize 2 --no-site